# airootfs/etc/systemd/system/compton.service
[Unit]
Description=Compton systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/compton

[Install]
WantedBy=multi-user.target
